# xeokey.service
[Unit]
Description=XeoKey Password Manager
Documentation=https://github.com/xeoxaz/XeoKey
After=network.target mongodb.service
Wants=mongodb.service

[Service]
Type=simple
User=xeo
Group=xeo
WorkingDirectory=/home/xeo/xeokey
Environment=NODE_ENV=production
Environment=SYSTEMD_SERVICE=true
EnvironmentFile=-/home/xeo/xeokey/src/.env
ExecStart=/usr/bin/bun run host
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xeokey

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/xeo/xeokey/logs /home/xeo/xeokey/src/backups
RemoveIPC=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Update handling - allow git operations
ReadWritePaths=/home/xeo/xeokey

[Install]
WantedBy=multi-user.target
